# kobo-dotfile-hack/scripts/98-dotfile.rules
# It should run close to on-animator, so onboard *should* be mounted
KERNEL=="loop0", RUN+="/usr/local/geek1011/bin/dotfile-patch.sh"
